# root/etc/s6-overlay/s6-rc.d/ttyd/run
#!/command/with-contenv /bin/bash
set -e

exec /etc/s6-overlay/scripts/ttyd.sh
